# ha-addon/rootfs/etc/services.d/spectrum_knx/run
#!/usr/bin/with-contenv bashio

# Wait for Postgres to be ready
bashio::log.info "Waiting for PostgreSQL to start..."
until s6-setuidgid postgres /usr/lib/postgresql/15/bin/pg_isready -h 127.0.0.1 -p 5432; do
    sleep 1
done
bashio::log.info "PostgreSQL is ready!"

# DB Configuration
export POSTGRES_HOST=127.0.0.1

# Read KNX options from Add-on config if provided
for var in KNX_PROJECT_PATH KNX_PASSWORD KNX_GATEWAY_IP KNX_GATEWAY_PORT KNX_CONNECTION_TYPE KNX_LOCAL_IP KNX_INDIVIDUAL_ADDRESS KNX_ROUTE_BACK KNX_MULTICAST_GROUP KNX_MULTICAST_PORT KNX_KNXKEYS_FILE KNX_KNXKEYS_PASSWORD KNX_SECURE_USER_ID KNX_SECURE_USER_PASSWORD KNX_SECURE_DEVICE_PASSWORD KNX_SECURE_BACKBONE_KEY KNX_SECURE_LATENCY_MS LOG_LEVEL; do
    if bashio::config.has_value "$var"; then
        export "$var"="$(bashio::config "$var")"
    fi
done

if [ -z "$LOG_LEVEL" ]; then
    export LOG_LEVEL="INFO"
fi

# Ingress support path (provided by Home Assistant)
if bashio::config.has_value 'ingress_port'; then
    # We serve directly on the port and Home Assistant handles the Ingress proxying
    export PORT=8000
fi

# Set root path if running under HA Ingress
# Home Assistant supervisor sets standard HTTP headers, but we might need to tell FastAPI
# about the root_path if we want to build URLs properly, but relative paths usually work.



cd /app/backend
bashio::log.info "Starting Spectrum KNX Backend..."
exec /opt/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
